# init/s6/inetd/run
#!/bin/sh

exec /usr/sbin/inetutils-inetd --environment --foreground

# vim: ft=sh
