# init/s6/smtpd/run
#!/bin/sh

exec /usr/sbin/smtpd -d -f /etc/smtpd.conf

# vim: ft=sh
