# init/s6/syslogd/run
#!/bin/sh

exec /usr/bin/s6-socklog -d 4

# vim: ft=sh
